# services/dnsmasq-private-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Trigger Protonet dnsmasq private restart

[Path]
PathModified=/etc/protonet/system/wifi/enabled

[Install]
WantedBy=multi-user.target
