# services/dokku-hostname-protonet.path
# ExperimentalPlatform
[Unit]
Description=dokku hostname changer

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target
