# services/hostname-avahi-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet hostname AVAHI

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target
