# services/hostname-ptw-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet Publish To Web Restarter

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target
