# services/hostname-protonet.path
[Unit]
Description=Protonet hostname changer

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target