# services/hostname-avahi-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet hostname AVAHI
After=network-online.target
Requires=network-online.target

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target
