# services/lock-password-protonet.path
# ExperimentalPlatform
[Unit]
Description=Experimental Platform Lock Account Trigger

[Path]
PathModified=/etc/protonet/system/ssh/lock

[Install]
WantedBy=multi-user.target
