# services/platform-haproxy-restart-on-gitlab-change.path
# ExperimentalPlatform
[Unit]
Description=Trigger a Restart of HAProxy if the gitlab IP address changes

[Path]
PathChanged=/etc/protonet/gitlab/ip

[Install]
WantedBy=multi-user.target
