# services/set-password-protonet.path
# ExperimentalPlatform
[Unit]
Description=Experimental Platform Passwor Setting Trigger

[Path]
PathModified=/etc/protonet/system/ssh/password

[Install]
WantedBy=multi-user.target
