# services/settimezone.service
[Unit]
Description=Set the timezone

[Service]
ExecStart=/usr/bin/timedatectl set-timezone Europe/Berlin
RemainAfterExit=yes
Type=oneshot
