# services/settimezone-protonet.service
# ExperimentalPlatform
[Unit]
Description=Set the timezone

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/timedatectl set-timezone Europe/Berlin
