# services/smb-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet SMB

[Path]
PathModified=/etc/protonet/hostname
PathModified=/data/samba/etc/

[Install]
WantedBy=multi-user.target
