# services/smb-restart-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet SMB

[Path]
PathModified=/etc/protonet/hostname

[Install]
WantedBy=multi-user.target
