# services/smb-shares-protonet.path
# ExperimentalPlatform
[Unit]
Description=Protonet SMB: Trigger shares reconfiguration

[Path]
PathModified=/etc/samba/

[Install]
WantedBy=multi-user.target
