# services/trigger-update-protonet.path
# ExperimentalPlatform
[Unit]
Description=Experimental Platform Channel Path Trigger
After=network-online.target
Requires=network-online.target

[Path]
PathModified=/etc/protonet/system/channel

[Install]
WantedBy=multi-user.target
