# services/trigger-update-protonet.path
# ExperimentalPlatform
[Unit]
Description=Experimental Platform Channel Path Trigger
After=dokku-protonet.service network-online.target
Requires=dokku-protonet.service network-online.target

[Path]
PathModified=/etc/protonet/system/channel

[Install]
WantedBy=multi-user.target
