# services/zfs-status.service
[Unit]
Description=Checks zpool status of protonet storage and writes it to skvs

[Service]
Type=oneshot
ExecStart=/opt/bin/zfs-status